# package/DEBIAN/postinst
#!/bin/sh
systemctl enable cluster-lab.service
