# package/DEBIAN/postinst
#!/bin/sh
systemctl enable hypriot-cluster.service
