# package/DEBIAN/postinst
#!/bin/sh
systemctl enable cluster-start.service
systemctl enable cluster-stop.service
